# app/deploy/logbook.service
# /etc/systemd/system/logbook.service
[Unit]
Description=Logbook workout tracker
After=network.target

[Service]
Type=simple
User=logbook
WorkingDirectory=/opt/logbook
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1
EnvironmentFile=/opt/logbook/.env
ExecStart=/usr/bin/node build
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
